# AccSevDetectionModel/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=python /app/backend/server.py
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/backend.out.log
stderr_logfile=/var/log/backend.err.log
environment=PYTHONUNBUFFERED=1

[program:frontend]
command=serve -s /app/frontend/dist -l 5173
autostart=true
autorestart=true
stdout_logfile=/var/log/frontend.out.log
stderr_logfile=/var/log/frontend.err.log